# deploy/btcedu-run.service
[Unit]
Description=btcedu - Process pending episodes through pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/AI-Startup-Lab/bitcoin-education
ExecStart=/home/pi/AI-Startup-Lab/bitcoin-education/.venv/bin/btcedu run-pending --max 3
EnvironmentFile=/home/pi/AI-Startup-Lab/bitcoin-education/.env
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
